# s3server/third_party/mero/scripts/install/usr/lib/systemd/system/m0t1fs@.service
[Unit]
Description=Mero client service
Requires=mero-kernel.service
After=mero-kernel.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/mero/mero-client start %i
ExecStop=/usr/libexec/mero/mero-client stop %i
RemainAfterExit=yes
KillMode=none
