# s3server/third_party/mero/scripts/install/usr/lib/systemd/system/mero-cleanup.service
[Unit]
Description=Mero cleanup

[Service]
Type=oneshot
ExecStart=/usr/libexec/mero/mero-cleanup
