# s3server/third_party/mero/scripts/install/usr/lib/systemd/system/mero-mkfs.service
[Unit]
Description=Mero mkfs

[Service]
Type=oneshot
Group=mero
ExecStart=/usr/libexec/mero/mero-mkfs
