# s3server/third_party/mero/scripts/install/usr/lib/systemd/system/mero-kernel.service
[Unit]
Description=Mero kernel module
Requires=lnet.service
Wants=mero-trace@kernel.service
After=lnet.service
Before=mero-trace@kernel.service mero-server@.service m0d@.service mero-server-confd.service mero-server-ha.service mero-client.service
PartOf=mero-server@.service m0d@.service mero-mkfs@.service mero-server-confd.service mero-server-ha.service mero-client.service m0t1fs@.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/mero/mero-kernel start
ExecStop=/usr/libexec/mero/mero-kernel stop
RemainAfterExit=yes
KillMode=none
TimeoutStopSec=1min
