# s3server/third_party/mero/scripts/install/usr/lib/systemd/system/mero-free-space-monitor.service
[Unit]
Description="Mero free space monitor"

[Service]
EnvironmentFile=/etc/sysconfig/mero-free-space-monitor
ExecStart=/usr/bin/python3 /usr/libexec/mero/mero-free-space-monitor $MERO_FREESPACE_ARGS
