# s3server/third_party/mero/scripts/install/usr/lib/systemd/system/mero.service
[Unit]
Description=Mero server helper
Requires=lnet.service
After=lnet.service
Before=mero-client.service
PartOf=mero-singlenode.service

[Service]
Type=oneshot
Group=mero
ExecStart=/usr/libexec/mero/mero start
ExecStartPost=-/usr/libexec/mero/mero-dixinit
ExecStop=/usr/libexec/mero/mero stop
RemainAfterExit=yes
KillMode=none
