# s3server/third_party/mero/scripts/provisioning/roles/iscsi-initiator/templates/iscsi-initiator-ssus.service.j2
[Unit]
Description=iSCSI initiator for SSUs target
Requires=iscsi-initiator-ssu1.service iscsi-initiator-ssu2.service
After=iscsi-initiator-ssu1.service iscsi-initiator-ssu2.service

[Service]
Type=oneshot
ExecStart=/bin/true
SyslogIdentifier=iscsi-initiator-ssus
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
